# rootfs/etc/s6-overlay/s6-rc.d/svc-gost/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

export GOST_LOGGER_LEVEL=${GOST_LOGGER_LEVEL:-info}

gost \
  -L "socks5://:${SOCKS5_PORT:-1080}" \
  -L "http://:${HTTP_PORT:-8080}"
